# app_2k1000la/systemd/xylt-safecloud-client.service.in
[Unit]
Description=XYLT SafeCloud evaluate polling client
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=@REPO_ROOT@
Environment=PYTHONUNBUFFERED=1
Environment=SAFECLOUD_DISCOVERY_PORT=8011
Environment=XYLT_SENSOR_SOURCE=mock
Environment=XYLT_2K0301_MQTT_HOST=127.0.0.1
Environment=XYLT_2K0301_MQTT_PORT=1883
Environment=XYLT_2K0301_MQTT_QOS=1
Environment=XYLT_2K0301_FIRST_TIMEOUT=8
Environment=XYLT_2K0301_STALE_AFTER=5
Environment=XYLT_SCENARIO_FILE=tests/scenarios/evaluate/gas_alarm.json
Environment=XYLT_OUTPUT_FILE=runtime/latest_evaluate_response.json
Environment=XYLT_TIMEOUT=30
Environment=XYLT_DISCOVERY_TIMEOUT=3
Environment=XYLT_INTERVAL=2
EnvironmentFile=-%h/.config/xylt/safecloud-client.env
ExecStart=/usr/bin/python3 @REPO_ROOT@/app_2k1000la/cloud_client.py --sensor-source ${XYLT_SENSOR_SOURCE} --scenario-file ${XYLT_SCENARIO_FILE} --output-file ${XYLT_OUTPUT_FILE} --timeout ${XYLT_TIMEOUT} --discovery-timeout ${XYLT_DISCOVERY_TIMEOUT} --include-debug --loop --interval ${XYLT_INTERVAL}
Restart=always
RestartSec=3
KillSignal=SIGINT

[Install]
WantedBy=default.target
