# deploy/systemd/2k0301/xylt-301-main.service.in
[Unit]
Description=XYLT 2K0301 sensor acquisition and actuator bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=@APP_DIR@
ExecStart=@APP_BIN@
Restart=always
RestartSec=3
KillSignal=SIGINT
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
